# 0microservices/gostreamcatcher/bin/goapi.service
[Unit]
Description=gostreamdownloader
After=multi-user.target

[Service]
User=ubuntu
Group=ubuntu
ExecStart=/bin/bash -c 'source /home/ubuntu/.bashrc && /home/ubuntu/godwn'

[Install]
WantedBy=multi-user.target